# deploy/supervisor/corun-worker.conf
[program:corun-worker]
command=/var/www/corun-ai/.venv/bin/python /var/www/corun-ai/worker.py --max-concurrent 1
directory=/var/www/corun-ai
user=admin
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=310
redirect_stderr=true
stdout_logfile=/var/log/corun-ai/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=HOME="/home/admin",LANG="C.UTF-8"
